# suse/x86_64/suse-leap-42.3-JeOS/root/usr/lib/systemd/system/data_part.service
[Unit]
Description=Add a data partition
Before=local-fs.target

[Service]
Type=oneshot
ExecStart=/bin/bash /sbin/create_data_part

[Install]
WantedBy=multi-user.target
